# munin/script/crontab
*/1 * * * * python /root/adt7410.py > /temporary/temp.log && chmod 777 /temporary/temp.log
*/720 * * * * bash /temporary/PWR/restart.sh
*/5 * * * * bash /root/speed.sh
